# production/systemd/irssh-panel.service
[Unit]
Description=IRSSH Panel Service
After=network.target postgresql.service redis.service
Requires=postgresql.service redis.service

[Service]
Type=simple
User=irssh
Group=irssh
WorkingDirectory=/opt/irssh-panel
Environment=NODE_ENV=production
Environment=PORT=3000
EnvironmentFile=/etc/irssh-panel/env
ExecStart=/usr/bin/node dist/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=irssh-panel
LimitNOFILE=65536

# Security
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
RestrictRealtime=true
RestrictNamespaces=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
